# docs/templates/autobot.service
[Unit]
Description=Autobot AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autobot
Group=autobot
WorkingDirectory=/var/lib/autobot
EnvironmentFile=/var/lib/autobot/.env
ExecStart=/usr/local/bin/autobot gateway

NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=true is included automatically when the workspace
# is outside of home directories (e.g. /var/lib/autobot).
ProtectHome=true
ReadWritePaths=/var/lib/autobot

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobot

[Install]
WantedBy=multi-user.target
